# deploy/systemd/dashfleet.service
[Unit]
Description=DashFleet system monitoring web UI
After=network.target

[Service]
Type=simple
User=dashboard
WorkingDirectory=/opt/dashfleet
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/dashfleet/.venv/bin/python -m main --web --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
